# api/docs/sample_files/lindele-web-watcher.path
[Unit]
Description=Monitors lindele web app for changes

[Path]
Unit=lindele-web-watcher.service
# Change this to the correct path for your system
PathChanged=/path/to/lindele/web

[Install]
WantedBy=multi-user.target